# traffic-control.service
[Unit]
Description=Aliyun CDT Traffic Control + Telegram Bot Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/aliyun-cdt-traffic-control
# 使用虚拟环境的 python 运行你最新的代码文件，不带多余的旧参数
ExecStart=/root/aliyun-cdt-traffic-control/venv/bin/python traffic_control.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
